# recipes-examples/example-systemd/files/example-systemd.service
[Unit]
Description=example-systemd

[Service]
ExecStart=/usr/bin/example_systemd.sh

[Install]
WantedBy=multi-user.target
